# bagian5/case4/platform/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:web]
directory=/var/www/
command=apachectl start
